# config.cron
0 0 * * * /home/michele/Documents/backup_linux.sh